# services/squid/run
#!/bin/sh
exec chpst -u squid:nogroup squid -f ./squid.conf -sN
