# services/privoxy.1/run
#!/bin/sh
exec 2>&1
exec chpst -u privoxy:nogroup privoxy --no-daemon
